# package/DEBIAN/postinst
#!/bin/bash
systemctl stop nginx
systemctl stop php8.1-fpm
systemctl stop systemd-journald

mv -n /config.json /etc/woodevops-toolkit/config.json

cp -f /etc/woodevops-toolkit/custom_files/php.ini /etc/php/8.1/fpm/php.ini
cp -f /etc/woodevops-toolkit/custom_files/nginx.conf /etc/nginx/nginx.conf
cp -f /etc/woodevops-toolkit/custom_files/journald.conf /etc/systemd/journald.conf

php /etc/woodevops-toolkit/scripts/update_php_ini
php /etc/woodevops-toolkit/scripts/update_journald_conf
php create_nginx_config_files all

systemctl enable nginx
systemctl start nginx
systemctl enable php8.1-fpm
systemctl start php8.1-fpm
systemctl enable systemd-journald
systemctl start systemd-journald