# package/DEBIAN/postinst
#!/bin/bash
systemctl stop nginx
systemctl stop php8.1-fpm
systemctl stop systemd-journald

mv -n /config.json /etc/woodevops-toolkit/config.json

php /etc/woodevops-toolkit/scripts/update_php_ini
php /etc/woodevops-toolkit/scripts/update_journald_conf
php create_nginx_config_files all

systemctl enable nginx
systemctl start nginx
systemctl enable php8.1-fpm
systemctl start php8.1-fpm
systemctl enable systemd-journald
systemctl start systemd-journald